# src/profile.d/xbuild-repos.sh
#!/usr/bin/bash

alias xbuild-repos='clear;xbuild-repos'
