# src/profile.d/xdeploy.sh
#!/usr/bin/bash

alias xdeploy='clear;xdeploy -v'
